# runservers.service
[Unit]
Description=script run

[Service]
ExecStart=/etc/run_script.sh
Type=oneshot
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
